# vnc-app/supervisor/app.conf
[program:qt-arm-py-ui]
environment=DISPLAY=:1
directory=/usr/src/app
command=python3 -u myuiapp.py
startretries=1
autorestart=false